# submission_debugger/deploy/submission_debugger.service
[Unit]
Description=Submission Debugger (FastAPI)
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/cvpr_competition/submission_debugger
EnvironmentFile=/opt/cvpr_competition/submission_debugger/.env
ExecStart=/opt/cvpr_competition/submission_debugger/scripts/run_foreground.sh
Restart=always
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=30
StandardOutput=append:/opt/cvpr_competition/submission_debugger/data/server.log
StandardError=append:/opt/cvpr_competition/submission_debugger/data/server.log

[Install]
WantedBy=multi-user.target
